# deploy/systemd/minsaeng100-collector.timer
[Unit]
Description=Minsaeng100 daily indicator collection timer

[Timer]
OnCalendar=*-*-* 06:20:00
RandomizedDelaySec=20m
Persistent=true
Unit=minsaeng100-collector.service

[Install]
WantedBy=timers.target
